# services/llm_factory/supervisord.conf
[supervisord]
nodaemon=true

[program:flaskapp]
command=python api.py
autostart=true
autorestart=true
